# deploy/hredh-backend.service
[Unit]
Description=Hidden Ridge EDH — FastAPI Backend
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/hiddenridgeedh/backend
EnvironmentFile=/var/www/hiddenridgeedh/backend/.env
ExecStart=/var/www/hiddenridgeedh/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8003 \
    --workers 2 \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hredh-backend

[Install]
WantedBy=multi-user.target
